# twemproxy/supervisord.conf
[supervisord]
nodaemon=true

[program:nutcracker]
command=nutcracker --conf-file=/etc/twemproxy/conf/twemproxy.yml