# twemproxy/supervisord.conf
[supervisord]
nodaemon=true

[program:nutcracker]
command=/bin/bash -c "nutcracker --conf-file=/etc/twemproxy/conf/twemproxy.yml"